# tools/codebox/codebox.conf
[program:codebox]
command=codebox run /var/workspace -p 8001 -e beetaa@outlook.com -u beetaa:beetaa
user=root
