# tools/zed/zedd.conf
[program:zedd]
command=zedd --user beetaa --pass beetaa --root /var/workspace --remote --enable-run
user=root
